# conf/healthCheck.ini
[program:healthCheck]
command=/usr/bin/python3 /etc/supervisord.d/scripts/healthCheck
process_name=%(program_name)s
numprocs=1
directory=/etc/supervisord.d/scripts/
autostart=false
startsecs=2
startretries=5
autorestart=true
stopsignal=TERM
stopwaitsecs=3
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=true
stdout_logfile=/etc/supervisord.d/scripts/healthCheck.log